# virt-archivos/udev-usbhotplug/10-mandoXbox.rules
ATTRS{idVendor}=="0e6f", ATTRS{idProduct}=="02a8", RUN+="/usr/local/bin/retrobox-mandos"
